# deploy/supervisor.conf
[program:shufa]
command=/var/www/shufa/venv/bin/gunicorn -c deploy/gunicorn.conf wsgi:app
directory=/var/www/shufa
user=root
autostart=true
autorestart=true
environment = MODE="PRODUCTION"